# orchestration/templates/task-server.conf
description "Task server"
author "Ian Duncan"

start on [2345]
stop on [016]

setuid task_server
chdir /var/www/tasks
exec task_server

respawn
respawn limit 10 60
